# config/com.peonnotify.watchdog-cron.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.peonnotify.watchdog-cron</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>$HOME/.claude/hooks/peon-watchdog.sh</string>
    <string>--cron</string>
  </array>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>StandardOutPath</key>
  <string>$HOME/.claude/logs/watchdog-cron.log</string>
  <key>StandardErrorPath</key>
  <string>$HOME/.claude/logs/watchdog-cron.log</string>
  <key>AbandonProcessGroup</key>
  <true/>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
